# window_manager/wallpaper.timer
[Unit]
Description=change wallpaper periodically

[Timer]
OnCalendar=*-*-* 04:00:00
Unit=wallpaper.service

[Install]
WantedBy=timers.target
