# window_manager/wallpaper.service
[Unit]
Description=set random wallpaper

[Service]
ExecStart=/usr/local/bin/wallpaper --random

[Install]
WantedBy=multi-user.target
