# desktop/systemd/trash.service
[Unit]
Description=clear old trash files

[Service]
Type=oneshot
ExecStart=trash-empty 15

[Install]
WantedBy=multi-user.target
