# desktop/systemd/wallpaper.timer
[Unit]
Description=change wallpaper periodically

[Timer]
OnCalendar=Sat *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
