# system/openweather.timer
[Unit]
Description=update weather information periodically

[Timer]
OnCalendar=*:0/6
Persistent=true

[Install]
WantedBy=timers.target
